# quarry-web-dev.service
[Unit]
Description=Development web server for quarry

[Service]
ExecStart=/srv/venv/bin/python /vagrant/quarry.wsgi
User=vagrant
group=vagrant
